# systemd/user/moon-update.timer
[Unit]
Description=Run Moon Wallpaper Update

[Timer]
OnBootSec=30s
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
