# systemd/user/moon-update.service
[Unit]
Description=Moon Wallpaper Update

[Service]
Type=oneshot
ExecStart=%h/.config/moonphase/moon_update.sh
